# cdn/supervisor.d/cdnapi.ini
[program:cdnapi]
# user=nginx
command=/usr/local/bin/cdn cdnapi
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true